# roles/basevm_create/files/regenerate-ssh-keys.service
[Unit]
Description=Generate SSH host keys
ConditionFirstBoot=yes
Before=ssh.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A
ExecStartPost=/usr/bin/systemctl disable regenerate-ssh-keys.service
ExecStartPost=/usr/bin/rm /etc/systemd/system/regenerate-ssh-keys.service

[Install]
WantedBy=multi-user.target
